# deploy/provision/roles/reverse_proxy/files/cert-renew.timer
[Unit]
Description=Certificate Renewal

[Timer]
Persistent=true
OnCalendar=*-*-* 5:38
RandomizeDelaySec=1h

[Install]
WantedBy=timers.target
